# digitemp-node-exporter.timer
[Unit]
Description=Timer to run prometheus digitemp exporter

[Timer]
OnCalendar=*:*:0/10
Persistent=true

[Install]
WantedBy=timers.target
